# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[program:gogogon-consumer]
command=consumer.py              ; the program (relative uses PATH, can take args)
autorestart=true            ; always restart when program exits, regardless of exit code
startretries=10

[include]
files=conf/supervisord-server.conf

[eventlistener:gogogon-listener]
command=listener.py
events=PROCESS_STATE_FATAL
